# supervisord.conf
[supervisord]
nodaemon=true


[program:geo]
command=python workflow.py
process_name=geo%(process_num)s
numprocs=4
autostart=true
autorestart=unexpected
redirect_stderr=true
stdout_logfile=geo.log
